# overlay/hive-os/hive-zorg/files/hive-zorg.service
[Unit]
Description=HIVE ZORG-Ω Agent - The Shadow (Security & Zero-Trust Enforcement)
Documentation=man:zorg-omega(8)
After=network-online.target hive-personality.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/libexec/hive/zorg-omega.py
Restart=always
RestartSec=5
User=hive
Group=hive

# High security hardening
NoNewPrivileges=true
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectClock=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

ReadWritePaths=/etc/hive /var/lib/hive /var/log/hive

# Capabilities for security operations
AmbientCapabilities=CAP_AUDIT_WRITE CAP_NET_RAW CAP_NET_ADMIN
Capabilities=cap_audit_write,cap_net_raw,cap_net_admin=ep

[Install]
WantedBy=multi-user.target
