# overlay/hive-os/hive-systemd/files/hive-collective-sync.timer
[Unit]
Description=HIVE Collective Synchronization Timer - Periodic agent sync
Requires=hive-collective-sync.service

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true
AccuracySec=10s

[Install]
WantedBy=timers.target
